# systemd/invisigoth.service
[Unit]
Description=Invisigoth AI Self-Bootstrapping Assistant
After=network.target

[Service]
Type=simple
User=invisigoth
WorkingDirectory=/opt/invisigoth
ExecStart=/opt/invisigoth/venv/bin/python3 /opt/invisigoth/main.py
Restart=on-failure
EnvironmentFile=-/opt/invisigoth/.env
StandardOutput=append:/opt/invisigoth/logs/app.log
StandardError=append:/opt/invisigoth/logs/app.log
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
